# Agendamento/nextcloudcron.service
[Unit]
Description=Nextcloud cron.php job

[Service]
ExecStart=docker exec -t -u www-data nextcloud php -f /var/www/html/cron.php
KillMode=process
